# layers/meta-tegrademo/recipes-demo/data-overlay-setup/data-overlay-setup/data-overlay-setup.service.in
[Unit]
Description=Create directories for overlay mounts
DefaultDependencies=no
RequiresMountsFor=/data

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@SBINDIR@/data-overlay-setup

[Install]
WantedBy=local-fs.target
